# systemd/moonshot-v2-social.timer
[Unit]
Description=Moonshot v2 — Social Collection 1h Timer

[Timer]
OnCalendar=*-*-* *:30:00
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
